# deploy/systemd/homescreen.service.example
[Unit]
Description=Homescreen App
After=network-online.target display-manager.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/homescreen
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/usr/bin/bash -lc 'for i in {1..30}; do [ -S /tmp/.X11-unix/X0 ] && exit 0; sleep 1; done; echo "X display :0 not ready"; exit 1'
ExecStart=/home/pi/homescreen/.venv/bin/python /home/pi/homescreen/main.py
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
